# scenarios/georgetown/files/mock-build.service
[Unit]
Description=Mock Build Application
After=network.target

[Service]
ExecStart=/home/admin/mock_build.sh

[Install]
WantedBy=multi-user.target
